# installer_scripts/deb/membrane-service-proxy/DEBIAN/postinst
#!/bin/bash
update-rc.d -f service_proxy remove
